# root/etc/services.d/devwallet/run
#!/usr/bin/with-contenv bash

cd /code || exit
exec \
	s6-setuidgid abc npm run serve > /applogs/app.log
